# .docker/php8.1-fpm-alpine/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/etc/supervisor/supervisord.log
user=_USER_

[program:php-fpm]
command=php-fpm

[program:nginx]
command=nginx

